# systemd/flow.service
[Unit]
Description=Flow Voice Dictation Daemon
After=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.local/bin/flow run
Restart=on-failure
RestartSec=5
Environment=RUST_LOG=info

[Install]
WantedBy=default.target
